# rootfs/etc/s6-overlay/s6-rc.d/init-grimmory-permissions/run
#!/bin/sh
set -e

# Ensure data, bookdrop, and books directories exist and are writable by the target user
mkdir -p /config/grimmory /bookdrop /books
chown grimmory:grimmory /config /bookdrop /books 2>/dev/null || true
