# rootfs/etc/s6-overlay/s6-rc.d/init-nginx-directories/run
#!/bin/sh
set -e

# Ensure nginx has a place for the pidfile
mkdir -p /var/run/nginx
